# debian/etc/systemd/user/gnome-xsession@.service.d/override.conf
[Service]
ExecStart=/etc/X11/Xsession gnome-session
