# common/etc/systemd/user/gnome-session-signal-init.service.d/override.conf
[Service]
ExecStart=bash -c 'echo Done >%T/gnome-session-wait-running-fifo@${DISPLAY}'
