# fedora/etc/systemd/system/xvfb@.service.d/override.conf
[Service]
Type=forking
ExecStart=
ExecStart=/usr/local/bin/xvfb.sh -ac -listen tcp %I -screen 0 "${XVFB_RESOLUTION}" -fbdir "${XVFB_FBDIR}"
