# fedora/etc/systemd/user/gnome-xsession@.service.d/override.conf
[Service]
ExecStart=/etc/X11/xinit/Xsession gnome-session
