# common/etc/systemd/user/org.gnome.Shell@wayland.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/gnome-shell --nested --wayland --no-x11
