# common/etc/systemd/user/dbus-proxy@.socket
[Unit]
Description=D-Bus TCP proxy on port %I

[Socket]
ListenStream=%I

[Install]
WantedBy=sockets.target
